# supervisord.jboss.conf
[program:jboss]
command=/usr/share/jboss/bin/standalone.sh -b 0.0.0.0
startsecs=30
user=jboss
